# docker/mirror/archlinux/crontab
16 * * * * root docker compose -p mirror exec archlinux bash /sync
